# overlay.d/20userlogin/usr/lib/systemd/system/silverblue-run-initial-setup.service
[Unit]
ConditionFirstBoot=true
Before=gdm.service

[Service]
ExecStart=/bin/touch /var/lib/gdm/run-initial-setup

[Install]
WantedBy=multi-user.target
