# overlay.d/20userlogin/usr/lib/systemd/system/silverblue-core-gdm-autologin.service
[Unit]
ConditionPathExists=/run/ostree-live
Before=gdm.service

[Service]
ExecStart=/usr/libexec/silverblue-core-gdm-autologin

[Install]
WantedBy=multi-user.target
